# templates/letschat.erb
#!/bin/sh
#
# /etc/init.d/letschat
# Subsystem file for letschat application server
#
# chkconfig: 2345 95 05	(1)
# description: letschat application server daemon
#
# processname: letschat
# config: <%= @deploy_dir %>/settings.yml
# config: /etc/sysconfig/letschat
# pidfile: /var/run/letschat.pid

# source function library
. /etc/rc.d/init.d/functions

# pull in sysconfig settings
[ -f /etc/sysconfig/letschat ] && . /etc/sysconfig/letschat

RETVAL=0
prog="letschat"

start() {
	echo -n $"Starting $prog:"
        (cd <%= @deploy_dir %> && node app.js) &
	RETVAL=$?
	[ "$RETVAL" = 0 ] && touch /var/lock/subsys/$prog
	echo
}

stop() {
	echo -n $"Stopping $prog:"
	killproc $prog -TERM
	RETVAL=$?
	[ "$RETVAL" = 0 ] && rm -f /var/lock/subsys/$prog
	echo
}

reload() {
	echo -n $"Reloading $prog:"
	killproc $prog -HUP
	RETVAL=$?
	echo
}

case "$1" in
	start)
		start
		;;
	stop)
		stop
		;;
	restart)
		stop
		start
		;;
	reload)
		reload
		;;
	condrestart)
		if [ -f /var/lock/subsys/$prog ] ; then
			stop
			# avoid race
			sleep 3
			start
		fi
		;;
	status)
		status $prog
		RETVAL=$?
		;;
	*)	(10)
		echo $"Usage: $0 {start|stop|restart|reload|condrestart|status}"
		RETVAL=1
esac
exit $RETVAL
